# com.payton.devvm-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.payton.devvm-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/paytonguo/Repos/devvm-setup/sync-devvms.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/paytonguo/.devvm-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/paytonguo/.devvm-sync.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/paytonguo</string>
    </dict>
</dict>
</plist>
